# ASL3-EZ.service
###############################################
###             ASL3-EZ Service             ###
###############################################
# https://www.github.com/GooseThings/ASL3-EZ/
#
[Unit]
Description=ASL3-EZ AllStarLink 3 Node Manager
After=network.target asterisk.service
Wants=asterisk.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/ASL3-EZ
ExecStart=/opt/ASL3-EZ/venv/bin/gunicorn --workers 1 --bind 0.0.0.0:5000 --timeout 60 --log-level info app:app
Restart=on-failure
RestartSec=5

# ── AMI Credentials ─────────────────────────────────────────────────────
# Set these to match the [username] stanza in /etc/asterisk/manager.conf
# The user must have: write = system,call,log,verbose,command,agent,user,config,dtmf,reporting,cdr,dialplan
Environment="AMI_USER=admin"
Environment="AMI_SECRET=yourpassword"

# ── Paths ────────────────────────────────────────────────────────────────
Environment=RPT_CONF_PATH=/etc/asterisk/rpt.conf
Environment=MANAGER_CONF=/etc/asterisk/manager.conf
Environment=BACKUP_DIR=/etc/asterisk/rpt_backups
Environment=DB_PATH=/etc/asterisk/asl3ez.db

# ── Network ──────────────────────────────────────────────────────────────
Environment=PORT=5000
Environment=HOST=0.0.0.0
Environment=AMI_HOST=127.0.0.1
Environment=AMI_PORT=5038

# ── Security ─────────────────────────────────────────────────────────────
# Change this to a long random string in production
Environment=SECRET_KEY=asl3-ez-change-me-in-production

# ── Logging ───────────────────────────────────────────────────────────────
# LOG_LEVEL controls ASL3-EZ's own log verbosity in journald.
# INFO (default) shows normal operation messages; DEBUG adds detailed trace
# output for every AMI command, config parse, and API call.
# Uncomment to enable verbose debug logging:
#Environment=LOG_LEVEL=DEBUG
# Path to the Asterisk log file shown in the Asterisk Console page.
Environment=ASTERISK_LOG_PATH=/var/log/asterisk/messages.log

[Install]
WantedBy=multi-user.target
